# files/puppetmaster.upstart.conf
description "PuppetMaster unicorn server"

start on runlevel [2345]
stop on runlevel [016]

chdir /etc/puppet/rack
exec unicorn -c /etc/puppet/rack/unicorn.rb /etc/puppet/rack/config.ru
